# BaseSystem/System/Library/LaunchDaemons/com.apple.lightsoutmanagementRecoveryOSd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>KeepAlive</key>
	<false/>
	<key>Label</key>
	<string>com.apple.lightsoutmanagementRecoveryOSd</string>
	<key>EnableTransactions</key>
	<false/>
	<key>POSIXSpawnType</key>
	<string>Interactive</string>
	<key>EnablePressuredExit</key>
	<false/>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/libexec/lightsoutmanagementRecoveryOSd</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>_LimitLoadFromBootMode</key>
	<string>fvunlock</string>
</dict>
</plist>
